# service/jack-the-ripper.service
[Service]
Type=oneshot
ExecStart=/opt/jack-the-ripper/jack-the-ripper.sh
