# systemd/jack-the-ripper@.service
[Service]
User=eric
Group=eric
Type=oneshot
ExecStart=/opt/jack-the-ripper/jack-the-ripper.sh %I
