# udev/51-jack-the-ripper.rules
ACTION=="change", KERNEL=="sr0", RUN+="/opt/jack-the-ripper/jack-the-ripper.sh"
